# crontab
0 2 * * * root /runbackup.sh >> /var/log/cron.log 2>&1

